# pi/todo-receipts-printer.service
[Unit]
Description=Todo Receipts Print Polling Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ashni
WorkingDirectory=/home/ashni/todo-receipts
Environment=CLOUD_URL=https://todo-receipts-production.up.railway.app
ExecStart=/usr/bin/node /home/ashni/todo-receipts/poll-and-print.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
